# monitor_test.service

[Unit]
Description=Monitor process "test"

[Service]
Type=oneshot
ExecStart=/usr/local/bin/monitor_test.sh
